# var/lib/invis-bbu/80-backupdisk.rules.template
SUBSYSTEMS=="usb", KERNEL=="sd*", ATTRS{serial}=="serialnr", ATTR{size}=="partsize", SYMLINK+="backup", RUN+="/usr/bin/udbdiskplugged"
