# usr/share/doc/packages/invis-bbu/server-parts/invisbu.cron
# Umgebungsvariablen setzen
PATH="/sbin:/usr/sbin:/usr/local/sbin:/root/bin:/usr/local/bin:/usr/bin:/bin:/usr/lib/mit/bin:/usr/lib/mit/sbin"
# Email Versand fuer die invis cron-jobs deaktiviert. Es sind einfach zu viele Mails
MAILTO=/dev/null

# Taeglich den verbleibenden Platz auf Volume /srv pruefen
0 5 * * * root /usr/bin/chkbuspace
